# utilities/duckdns.timer
[Unit]
Description=DuckDNS Update Every Four Hours

[Timer]
OnCalendar=00/4:00
Persistent=true
Unit=duckdns.service

[Install]
WantedBy=timers.target