# services/photoprism/photoprism-index.service
[Unit]
Description=PhotoPrism Index

[Service]
ExecStart=/usr/bin/podman exec -t photoprism photoprism index -c
Restart=no
Type=oneshot